# update-motd.d/90-footer
#!/bin/sh
#
# 90-footer - write the admin's footer to the MOTD
#

[ -f /etc/motd.tail ] && cat /etc/motd.tail || true
